# docker/xfce4/dev/etc/supervisor/conf.d/update-hosts.conf
[program:updatehosts]
command=bash -c 'sleep 5 && /usr/bin/update-hosts.sh'
user=root
autostart=true
autorestart=false
priority=900
